# sys-apps/systemd-rmilter/files/rmilter.service
[Service]
ExecStart=/usr/sbin/rmilter -c /etc/rmilter/rmilter.conf
User=rmilter

[Install]
WantedBy=multi-user.target
