# app-misc/n8n/files/n8n.service
[Unit]
Description = n8n service
After = network-online.target

[Service]
User = n8n
Restart = on-failure
RestartSec = 5
EnvironmentFile = -/etc/env.d/77n8n
ExecStart = /usr/bin/n8n

[Install]
WantedBy = multi-user.target
