# net-mail/mailman/files/mailman3-cron.timer
[Unit]
Description=Run mailman Cronjobs

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target