# net-mail/mailman/files/mailman3.cron
#!/bin/sh
/sbin/runuser -u mailman -- nice -n10 /usr/bin/mailman digests --periodic
/sbin/runuser -u mailman -- nice -n10 /usr/bin/mailman notify