# app-backup/AutoMySQLBackup/files/automysqlbackup.service
[Unit]
Description=Auto Mysql Backup

[Service]
ExecStart=/usr/bin/automysqlbackup
User=root