# app-backup/AutoMySQLBackup/files/automysqlbackup.timer
[Unit]
Description=Auto Mysql Backup

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
