# sys-apps/systemd-custom-timer/files/clamav-unofficial-sigs.service
[Unit]
Description=ClamAV Unofficial Signatures Updater

[Service]
Type=oneshot
ExecStart=/usr/sbin/clamav-unofficial-sigs.sh
