# sys-apps/systemd-custom-timer/files/geoipupdate.service
[Unit]
Description=GeoIP Database Updater

[Service]
Type=oneshot
ExecStart=/usr/sbin/geoipupdate.sh --force
