# sys-apps/systemd-custom-timer/files/geoipupdate.timer
[Unit]
Description=timer for geoipupdate.service

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
