# sys-apps/systemd-custom-timer/files/sa-update.timer
[Unit]
Description=timer for sa-update.service

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
