# sys-apps/systemd-custom-timer/files/start-failed.service
[Unit]
Description=try start failed services

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl default
RemainAfterExit=yes
