# sys-apps/systemd-custom-timer/files/start-failed.service
[Unit]
Description=try start failed services

[Service]
ExecStart=/bin/systemctl default
