# sys-apps/systemd-custom-timer/files/start-failed.timer
[Unit]
Description=timer for start-failed.service 

[Timer]
OnBootSec=1min

[Install]
WantedBy=timers.target
