# sys-apps/systemd-dovecot-timer/files/doveadm-quota-recalc.timer
[Unit]
Description=doveadm quota recalc

[Timer]
OnCalendar=*-*-* 00:45:00

[Install]
WantedBy=timers.target
