# sys-apps/systemd-dovecot-timer/files/doveadm-force-resync.timer
[Unit]
Description=doveadm force-resync

[Timer]
OnCalendar=*-*-* 04:30:00

[Install]
WantedBy=timers.target
