# sys-apps/systemd-dovecot-timer/files/doveadm-fts.timer
[Unit]
Description=doveadm fts

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
