# sys-apps/systemd-dovecot-timer/files/doveadm-replicate.timer
[Unit]
Description=doveadm replicate

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
