# sys-apps/systemd-dovecot-timer/files/doveadm-quota-recalc.service
[Unit]
Description=doveadm quota recalc

[Service]
Type=oneshot
ExecStart=/usr/bin/doveadm quota recalc -A