# sys-apps/systemd-dovecot-timer/files/doveadm-purge.timer
[Unit]
Description=doveadm purge

[Timer]
OnCalendar=*-*-* 04:15:00

[Install]
WantedBy=timers.target
