# sys-apps/systemd-dovecot-timer/files/doveadm-force-resync.service
[Unit]
Description=doveadm force-resync

[Service]
Type=oneshot
ExecStart=/usr/bin/doveadm force-resync -A '*'