# sys-apps/systemd-dovecot-timer/files/doveadm-replicate.timer
[Unit]
Description=doveadm replicate

[Timer]
OnCalendar=*-*-* 04:00:00

[Install]
WantedBy=timers.target
