# sys-apps/systemd-dovecot-timer/files/doveadm-purge.service
[Unit]
Description=doveadm purge

[Service]
Type=oneshot
ExecStart=/usr/bin/doveadm purge -A