# sys-apps/systemd-dovecot-timer/files/doveadm-fts.service
[Unit]
Description=doveadm fts

[Service]
Type=oneshot
ExecStart=/usr/bin/doveadm fts rescan -A
