# sys-apps/systemd-dovecot-timer/files/doveadm-quota-recalc.timer
[Unit]
Description=doveadm quota recalc

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
