# sys-apps/systemd-dovecot-timer/files/doveadm-purge.timer
[Unit]
Description=doveadm purge

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
